# scripts/systemd/mqtt-cert-renew.timer
[Unit]
Description=Run MQTT certificate renewal periodically

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=30m
Persistent=true
Unit=mqtt-cert-renew.service

[Install]
WantedBy=timers.target
